# contrib/systemd/aqs.service
[Unit]
Description=ember (aqs) Quickshell desktop shell
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/quickshell -p /usr/local/share/aqs/qml/shell.qml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
